# systemd/user/compositor.service
[Unit]
Description=Compositor
After=xorg.service
Wants=xorg.service

[Service]
ExecStart=xcompmgr -n

[Install]
WantedBy=wm.target
