# systemd/user/dark-mode.service
[Unit]
Description=Dark Mode
After=xsettingsd.service

[Service]
ExecStart=%h/Workspace/linux-utils/auto-gtk-theme/target/release/auto-gtk-theme
Restart=on-failure

[Install]
WantedBy=default.target
