# systemd/user/dark-mode.service
[Unit]
Description=Dark Mode
After=xsettingsd.service
Requires=xsettingsd.service

[Service]
Environment=RUST_LOG=debug
ExecCondition=dump_xsettings
ExecStart=%h/Workspace/linux-utils/night-mode/target/release/night-mode
Restart=on-failure
RestartSec=1

[Install]
WantedBy=default.target
