# systemd/user/dark-mode.service
[Unit]
Description=Dark Mode
Wants=xorg.service
After=xorg.service
Wants=systemd-networkd-wait-online.service
After=systemd-networkd-wait-online.service

[Service]
Environment=RUST_LOG=debug:night_mode
ExecCondition=i3 --get-socketpath
ExecStart=%h/workspace/linux-utils/night-mode/target/release/night-mode
Restart=always

[Install]
WantedBy=wm.target
