# systemd/user/devmon.service
[Unit]
Description=Mount plugged-in USB devices
After=default.target

[Service]
ExecStart=devmon
Restart=always

[Install]
WantedBy=default.target
