# systemd/user/event-reminder.timer
[Unit]
Description=Event Reminder

[Timer]
OnBootSec=1min
OnUnitActiveSec=20min

[Install]
WantedBy=timers.target
