# systemd/user/mailwatch@.service
[Unit]
Description=Email monitor: %i
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=%h/.dotfiles/bin/mailwatch "%i"
Restart=always
RestartSec=1

[Install]
WantedBy=wm.target
