# systemd/user/mbsync.timer
[Unit]
Description=Mailbox synchronization timer

[Timer]
OnBootSec=5s
OnUnitActiveSec=60s
Unit=mbsync.service

[Install]
WantedBy=timers.target
