# systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service

[Service]
Type=oneshot
ExecStart=fish -C 'mbsync -Va -c ~/.config/mbsync/config || true'
ExecStartPost=/usr/bin/fish /home/greg/.secconf/mbsync/postsync.sh

[Install]
WantedBy=wm.target
