# systemd/user/mbsync.service
[Unit]
Description=Mailbox synchronization service
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=mbsync -Va --pull -c %h/.secconf/mbsync/config
ExecStartPost=%h/.secconf/mbsync/postsync.sh

[Install]
WantedBy=wm.target
