# systemd/user/monitors.service
[Unit]
Description=Setup monitors
Wants=xorg.service
After=xorg.service


[Service]
ExecStart=enact -p top
#Type=oneshot
Restart=on-failure

[Install]
WantedBy=wm.target
