# systemd/user/monitors.service
[Unit]
Description=Setup monitors
After=xorg.service
Wants=xorg.service

[Service]
ExecStart=enact -p top
Type=oneshot
Restart=on-failure
RestartSec=0.1

[Install]
WantedBy=wm.target
