# systemd/user/polybar.service
[Unit]
Description=Polybar
After=monitors.service
Requires=monitors.service

[Service]
ExecStartPre=sleep 1
ExecStart=polybar -l trace primary
RestartSec=1
Restart=always

[Install]
WantedBy=wm.target
