# systemd/user/polybar.service
[Unit]
Description=Polybar
After=window-manager.service

[Service]
ExecStartPre=i3 --get-socketpath
ExecStart=polybar primary
RestartSec=1
Restart=always

[Install]
WantedBy=wm.target
