# systemd/user/polybar.service
[Unit]
Description=Polybar
After=monitors.service
Requires=monitors.service
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=sleep 1
ExecStart=polybar -l info primary
RestartSec=1
Restart=always

[Install]
WantedBy=wm.target
