# systemd/user/polybar.service
[Unit]
Description=Polybar
After=window-manager.service
Requires=window-manager.service

[Service]
ExecStart=polybar primary
Restart=always

[Install]
WantedBy=wm.target
