# systemd/user/redshift.service
[Unit]
Description=Blue light filter
PartOf=graphical-session.target

[Service]
ExecStart=redshift -c %h/.config/redshift.conf

[Install]
WantedBy=default.target
