# systemd/user/wallpaper.service
[Unit]
Description=Set wallpapers

[Service]
ExecStart=/home/greg/.dotfiles/bin/wallpaper

[Install]
WantedBy=wm.target
