# systemd/user/wallpaper.service
[Unit]
Description=Set wallpapers
After=monitors.service
Requires=monitors.service

[Service]
ExecStart=feh --no-fehbg --bg-scale %h/.config/wallpaper %h/.config/wallpaper

[Install]
WantedBy=wm.target
