# systemd/user/xresources.service
[Unit]
Description=Load Xresources
After=xorg.service

[Service]
ExecStart=xrdb %h/.dotfiles/x11/Xresources

[Install]
WantedBy=wm.target
