# systemd/user/xresources.service
[Unit]
Description=Load Xresources
After=xorg.service

[Service]
ExecStart=xrdb %h/.config/x11/xresources

[Install]
WantedBy=wm.target
